# backend/ops/systemd/mcu-mqtt-broker.service
[Unit]
Description=MCU Local MQTT Broker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__RUN_USER__
WorkingDirectory=__BACKEND_DIR__
ExecStart=/usr/bin/node scripts/mqtt-broker.mjs
Restart=always
RestartSec=3
Environment=NODE_ENV=production
Environment=MQTT_PORT=1883

[Install]
WantedBy=multi-user.target
